# mine.cron
1 4 * * * /usr/bin/curl -i http://localhost:8080/topurls >/tmp/topurls_cache.log 2>&1
11 4 * * * /usr/bin/curl -i http://localhost:8080/topdomains >/tmp/topdomains_cache.log 2>&1
* * * * * echo "Hello world" >> /var/log/cron.log 2>&1
